# packaging/systemd/lamco-rdp-server.service
# Lamco RDP Server - Systemd User Service
#
# This service runs the RDP server as a user service, using the user's
# Wayland session via XDG Portals.
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp lamco-rdp-server.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable lamco-rdp-server
#   systemctl --user start lamco-rdp-server
#
# For linger (start on boot without login):
#   sudo loginctl enable-linger $USER

[Unit]
Description=Lamco RDP Server
Documentation=https://github.com/lamco-admin/lamco-rdp-server
After=graphical-session.target
Wants=graphical-session.target

# Require a graphical session - we need Wayland compositor access
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=dbus
BusName=io.lamco.RdpServer

# Main executable with D-Bus service mode
ExecStart=/usr/bin/lamco-rdp-server --dbus-service

# Restart on failure, but not too aggressively
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Environment
Environment=RUST_LOG=info

# Filesystem hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectProc=invisible
ProcSubset=pid

# Allow writes to XDG runtime dir (Wayland socket, PipeWire, D-Bus)
ReadWritePaths=%t
# Allow writes to config and state directories
ReadWritePaths=%h/.config/lamco-rdp-server
ReadWritePaths=%h/.local/share/lamco-rdp-server

# Kernel and device hardening
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

# Restrict namespace creation
RestrictNamespaces=yes

# Network: only TCP (RDP listening)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Capability bounding (user service has no caps by default, but be explicit)
CapabilityBoundingSet=
AmbientCapabilities=

# Syscall filtering: allow only what an RDP server needs
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallArchitectures=native

[Install]
WantedBy=graphical-session.target
