# packaging/systemd/lamco-rdp-server-system.service
# Lamco RDP Server - Systemd System Service
#
# This service runs the RDP server as a system service. It requires
# a dedicated user account with access to a graphical session.
#
# NOTE: System-level RDP service requires special setup:
# - Create a dedicated lamco-rdp user
# - Configure session access via logind seat
# - Consider security implications of system-wide remote access
#
# For most use cases, the user service (lamco-rdp-server.service) is preferred.
#
# Installation:
#   sudo cp lamco-rdp-server-system.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable lamco-rdp-server-system
#   sudo systemctl start lamco-rdp-server-system

[Unit]
Description=Lamco RDP Server (System)
Documentation=https://github.com/lamco-admin/lamco-rdp-server
After=graphical.target
Wants=graphical.target

[Service]
Type=dbus
BusName=io.lamco.RdpServer.System

# Run as dedicated user (create with: sudo useradd -r -s /sbin/nologin lamco-rdp)
User=lamco-rdp
Group=lamco-rdp

# Main executable with D-Bus service mode on system bus
ExecStart=/usr/bin/lamco-rdp-server --dbus-service --system

# Restart on failure
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment
Environment=RUST_LOG=info
Environment=XDG_RUNTIME_DIR=/run/user/%U

# Filesystem hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectProc=invisible
ProcSubset=pid

# State directories (systemd creates these with correct ownership)
ConfigurationDirectory=lamco-rdp-server
StateDirectory=lamco-rdp-server
RuntimeDirectory=lamco-rdp-server

# Kernel and device hardening
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

# Restrict namespace creation
RestrictNamespaces=yes

# Network: only TCP (RDP listening) and Unix (D-Bus, Wayland)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Drop all capabilities - RDP server doesn't need any
CapabilityBoundingSet=
AmbientCapabilities=

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
